# files/XiGuaPi-V3/etc/init.d/wifi-reset
#!/bin/sh /etc/rc.common
# Copyright (C) 2022 MDSDTech <daniel@mdsdtech.my.id>

START=09   # runs very early, before PCIe modules load

boot() {
    # Assert reset (high) to make sure card is held in reset
    echo 113 > /sys/class/gpio/export
    echo out > /sys/class/gpio/gpio113/direction
    echo 1 > /sys/class/gpio/gpio113/value
    sleep 0.2

    # Deassert reset (low) to release card
    echo 0 > /sys/class/gpio/gpio113/value
    sleep 0.2

    # Force a PCIe rescan so the now-ready card is found
    echo 1 > /sys/bus/pci/rescan
}